# net-misc/rtp2httpd/files/rtp2httpd.init
#!/sbin/runscript

depend() {
	need net
}

description="Multicast RTP to HTTP translator"
command=/usr/bin/rtp2httpd
command_background=yes
pidfile=/var/run/rtp2httpd.pid
start_stop_daemon_args="-u nobody:nobody"
command_args="-q"
